# backend/scripts/crontab
# ==============================================================================
# GeoViable — Cron Schedule for Layer Updates
# ==============================================================================
# Runs on the 1st of every month at 03:00 UTC.
# Output is appended to /var/log/geoviable_update.log.
# ==============================================================================
0 3 1 * * cd /app && python -m scripts.update_layers >> /var/log/geoviable_update.log 2>&1
